# integration/exit.crontab
* * * * * * * echo "$MSG_START" && sleep 2 && echo "$MSG_DONE"
