# integration/override.crontab
VAR="hello from bar"
* * * * * * * echo "$VAR"
