# integration/zombie.crontab
* * * * * * * /bin/sleep 1 & exec /bin/sleep 0